# data/export/upstart/process.conf.erb
start on starting <%= app %>-<%= process.name %>
stop on stopping <%= app %>-<%= process.name %>
respawn

exec su - <%= user %> -c 'cd <%= engine.directory %>; export PORT=<%= port %>; export APP_NUM=<%= num %>; <%= process.command %> >> <%= log_root %>/<%=process.name%>-<%=num%>.log 2>&1'
