# data/export/upstart/process.conf.erb
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>
respawn

env PORT=<%= port %>
<% engine.env.each_pair do |var,env| -%>
env <%= var.upcase %>=<%= env %>
<% end -%>

exec su <%= user %> -l -m -c 'cd <%= engine.root %>; <%= process.command %> >> <%= log %>/<%=name%>-<%=num%>.log 2>&1'
