# data/export/upstart/process.conf.erb
start on starting <%= app %>-<%= process.name %>
stop on stopping <%= app %>-<%= process.name %>
respawn

exec su -s /bin/bash -c 'cd <%= engine.directory %>; export PORT=<%= port %>;<% engine.environment.each_pair do |var,env| %> export <%= var.upcase %>=<%= env %>; <% end %> exec "$0" "$@"' rails -- <%= process.command %> >> <%= log_root %>/<%= process.name %>-<%= num %>.log 2>&1
