# spec/resources/export/systemd/standard/app-bravo.target
[Unit]
StopWhenUnneeded=true
Wants=app-bravo-1.service
