# spec/resources/export/systemd/standard/app-alpha.target
[Unit]
StopWhenUnneeded=true
Wants=app-alpha-1.service
