# spec/resources/export/systemd/standard/app-alpha.target
[Unit]
PartOf=app.target
Wants=app-alpha@5000.service
