# spec/resources/export/systemd/standard/app-bravo.target
[Unit]
PartOf=app.target
Wants=app-bravo@5100.service
