# spec/resources/export/systemd/standard/app.target
[Unit]
Wants=app-alpha.target app-bravo.target app-foo_bar.target app-foo-bar.target

[Install]
WantedBy=multi-user.target
