# spec/resources/export/upstart/app-alpha-1.conf
start on starting app-alpha
stop on stopping app-alpha
respawn

exec su -s /bin/bash -c 'cd /tmp/app; export PORT=5000; exec "$0" "$@"' app -- ./alpha >> /var/log/app/alpha-1.log 2>&1
