# spec/resources/export/upstart/app.conf
pre-start script

bash << "EOF"
  mkdir -p /var/log/app
  chown -R app /var/log/app
EOF

end script

start on (started network-interface
          or started network-manager
          or started networking)

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking)
