# spec/resources/export/upstart/foreman-bravo-1.conf
start on starting foreman-bravo
stop on stopping foreman-bravo
respawn

chdir /Users/david/Code/foreman
exec su foreman -c 'export PORT=5100; ./bravo >> /var/log/foreman/bravo-1.log 2>&1'
