# spec/resources/export/upstart/foreman-alpha-2.conf
start on starting foreman-alpha
stop on stopping foreman-alpha
respawn

chdir /Users/david/Code/foreman
exec su foreman -c 'export PORT=5001; ./alpha >> /var/log/foreman/alpha-2.log 2>&1'
