# spec/resources/export/upstart/foreman-alpha.conf
start on starting foreman
stop on stopping foreman
