# spec/resources/export/systemd/concurrency/app-alpha.target
[Unit]
StopWhenUnneeded=true
Wants=app-alpha-1.service app-alpha-2.service
