# spec/resources/export/systemd/concurrency/app.target
[Unit]
StopWhenUnneeded=true
Wants=app-alpha.target

[Install]
WantedBy=multi-user.target
