# spec/resources/export/systemd/app-alpha.target
[Unit]
StopWhenUnneeded=true

[Install]
WantedBy=app.target
