# data/export/systemd/master.target.erb
[Unit]
StopWhenUnneeded=true
Wants=<%= process_master_names.join(' ') %>

[Install]
WantedBy=multi-user.target
