# data/export/systemd/process_master.target.erb
[Unit]
PartOf=<%= app %>.target
Wants=<%= process_names.join(' ') %>
