# data/export/systemd/process.service.erb
[Unit]
PartOf=<%= app %>-<%= name %>.target

[Service]
User=<%= user %>
WorkingDirectory=<%= engine.root %>
Environment=PORT=%i
<% if !engine.env.empty? -%>Environment=<% engine.env.each_pair do |var,env| %>"<%= var.upcase %>=<%= env %>" <% end %><% end -%>
ExecStart=/bin/bash -lc '<%= process.command %>'
Restart=always
StandardInput=null
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n
KillMode=mixed
TimeoutStopSec=<%= engine.options[:timeout] %>
